# myyd/build/target/supervisor/node/acl_node.ini
[program:acl_node]
user=root
command=/virus/cloud_waf_detect/acl/acl_node -mlog_level info &
directory=/virus/cloud_waf_detect/acl

stopsignal=INT
startsecs=3
stopasgroup=true
killasgroup=true

autostart=true
autorestart=true
stdout_logfile=/virus/supervisor/log/supervisor.log
stderr_logfile=/virus/supervisor/log/supervisor_err.log
[supervisord]