# myyd/build/target/supervisor/node/link_node.ini
[program:link_node]
user=root
command=/virus/cloud_waf_detect/link/link_node -mlog_level info &
directory=/virus/cloud_waf_detect/link

stopsignal=INT
startsecs=3
stopasgroup=true
killasgroup=true

autostart=true
autorestart=true
stdout_logfile=/virus/supervisor/log/supervisor.log
stderr_logfile=/virus/supervisor/log/supervisor_err.log
[supervisord]
