# myyy/build/target/supervisor/node/eye_node.ini
[program:eye_node]
user=root
command=/virus/cloud_waf_detect/eye/eye_node -mlog_level info &
directory=/virus/cloud_waf_detect/eye

stopsignal=INT
startsecs=3
stopasgroup=true
killasgroup=true

autostart=true
autorestart=true
stdout_logfile=/virus/supervisor/log/supervisor.log
stderr_logfile=/virus/supervisor/log/supervisor_err.log
[supervisord]
