# myyd/build/target/supervisor/server/link_server.ini
[program:link_server]
user=root
command=/virus/cloud_waf_detect/linkserver/link_server -mlog_level info &
directory=/virus/cloud_waf_detect/linkserver

stopsignal=INT
startsecs=3
stopasgroup=true
killasgroup=true

autostart=true
autorestart=true
stdout_logfile=/virus/supervisor/log/supervisor.log
stderr_logfile=/virus/supervisor/log/supervisor_err.log
[supervisord]
