# myyd/build/target/supervisor/server/saasyd_website.ini
[program:saasyd_website]
user=root
command=/virus/cloud_waf_detect/saas_services/saasyd_website/saasyd_website -mlog_level info &
directory=/virus/cloud_waf_detect/saas_services/saasyd_website

stopsignal=INT
startsecs=3
stopasgroup=true
killasgroup=true

autostart=true
autorestart=true
stdout_logfile=/virus/supervisor/log/supervisor.log
stderr_logfile=/virus/supervisor/log/supervisor_err.log
[supervisord]
