# myexample/build/supervisor/server/myexample.ini
[program:myexample]
user=root
command=/data/sa/service/myexample/myexample &
directory=/data/sa/service/myexample

stopsignal=INT
startsecs=3
stopasgroup=true
killasgroup=true

autostart=true
autorestart=true
stdout_logfile=/data/supervisor/log/supervisor.log
stderr_logfile=/data/supervisor/log/supervisor_err.log
[supervisord]